# templates/default/etc/profile.d/mysql_server.sh.erb
###
# Generated by Chef for [<%= node[:fqdn] %>]
# Environment: <%= node.chef_environment %>
###

alias <%= @calias %>='mysql --defaults-extra-file=<%= @cpath %>/etc/my.cnf -S <%= @cpath %>/var/mysqld.sock'
