# templates/default/etc/cron.d/binlog_cleaner.erb
###
# Generated by Chef for [<%= node[:fqdn] %>]
# Environment: <%= node.chef_environment %>
###

# Purges old binary log files if disk is getting full

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin

*/5 *	* * *	root	/usr/local/bin/binlogCleaner.sh > /dev/null 2>&1
